# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s
autorestart=true
priority=10

[program:xfce4]
command=startxfce4
environment=DISPLAY=":1",HOME="/home/ctfuser"
user=ctfuser
autorestart=true
priority=20

[program:container_api]
command=/opt/container_api/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8888
directory=/opt/container_api
environment=DISPLAY=":1"
autorestart=true
priority=40

[program:tmux_session]
command=bash -c "sleep 2 && tmux -S /tmp/ctf-tmux kill-session -t ctf 2>/dev/null; tmux -S /tmp/ctf-tmux new-session -d -s ctf -x 200 -y 50 && chmod 777 /tmp/ctf-tmux && while tmux -S /tmp/ctf-tmux has-session -t ctf 2>/dev/null; do sleep 10; done"
user=ctfuser
environment=HOME="/home/ctfuser"
autorestart=true
priority=45

[program:visible_terminal]
command=bash -c "sleep 4 && xfce4-terminal --maximize --title='CTF Agent Terminal' -e 'tmux -S /tmp/ctf-tmux attach-session -t ctf'"
user=ctfuser
environment=DISPLAY=":1",HOME="/home/ctfuser"
autorestart=true
priority=50
